# packaging/ndnsf-di-systemd/units/ndnsf-di-repo@.service
[Unit]
Description=NDNSF authoritative Repo %i
After=network-online.target nfd.service
Requires=nfd.service
StartLimitIntervalSec=60s
StartLimitBurst=3

[Service]
Type=simple
User=ndnsf-di-repo
Group=ndnsf-di-repo
EnvironmentFile=-/etc/ndnsf-di/repo-%i.env
ExecStart=/opt/ndnsf-di/current/bin/ndn-repo-ng ${NDNSF_DI_REPO_ARGS}
Restart=on-failure
RestartSec=5s
TimeoutStartSec=60s
TimeoutStopSec=30s
KillSignal=SIGTERM
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
ReadWritePaths=/var/lib/ndnsf-repo /run/ndnsf-di /var/log/ndnsf-di

[Install]
WantedBy=multi-user.target
