# packaging/ndnsf-di-systemd/units/ndnsf-di-provider@.service
[Unit]
Description=NDNSF-DI provider %i
After=network-online.target nfd.service ndnsf-di-controller.service
Requires=nfd.service ndnsf-di-controller.service
PartOf=ndnsf-di-providers.target

[Service]
Type=simple
User=ndnsf-di
Group=ndnsf-di
EnvironmentFile=-/etc/ndnsf-di/provider-%i.env
ExecStart=/opt/ndnsf-di/current/bin/di-native-provider ${NDNSF_DI_PROVIDER_ARGS}
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=60s
StartLimitBurst=3
TimeoutStartSec=120s
TimeoutStopSec=30s
KillSignal=SIGTERM
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
ReadWritePaths=/run/ndnsf-di /var/lib/ndnsf-di/cache /var/log/ndnsf-di

[Install]
WantedBy=ndnsf-di-providers.target
