# packaging/ndnsf-di-systemd/staging/nfd.service
[Unit]
Description=Isolated staging NFD dependency stub

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes
